# app/qzap/common/tool/check_submit.cron
# Example of cron to check submit
*/1 *  *  *  *  cd /data1/home/phongchen/work/project/gdt && date >> check_submit.log && /usr/local/bin/python app/qzap/common/tool/check_submit.py
